# turncableworker.service
[Unit]
Description=Turncable
After=network.target

[Service]
ExecStart=/home/pi/turncable/worker.sh
User=pi

[Install]
WantedBy=multi-user.target
